# packer/poweroff.timer
[Unit]
Description=Delayed poweroff

[Timer]
OnActiveSec=5
Unit=poweroff.target
